# old-arch/conf/deploy/gunicorn-csdndata.com.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/tmp/csdndata.com.socket

[Install]
WantedBy=sockets.target
